# systemd_units/shetserv.service
[Unit]
Description=SHET Server

[Service]
User=jobs
ExecStart=/home/jobs/.local/bin/shetserv

[Install]
WantedBy=multi-user.target
